# file_upload.service
[Unit]
Description=File Upload Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=a214
Group=a214
WorkingDirectory=/home/a214/VideoProject
ExecStart=/usr/bin/python3 file_upload_server.py
Restart=on-failure
RestartSec=5s
Environment=UPLOAD_FOLDER=/home/a214/result
Environment=PORT=8094
Environment=ALLOWED_EXTENSIONS=zip,rar,7z,tar.gz,tgz,txt,pdf,jpg,jpeg,png,gif,docx,xlsx,pptx
Environment=MAX_CONTENT_LENGTH=104857600

[Install]
WantedBy=multi-user.target